# packages/lakka/lakka_tools/usb-modeswitch-lakka/udev.d/90-rtl8812bu.rules
# RTL8812BU / RTL8811CU
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="0bda", ATTRS{idProduct}=="1a2b", RUN+="/usr/sbin/usb_modeswitch -K -v 0bda -p 1a2b"
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="0e8d", ATTRS{idProduct}=="2870", RUN+="/usr/sbin/usb_modeswitch -K -v 0e8d -p 2870"
